# codingclubupdate.conf
description "coding club update"
start on startup
task
exec /usr/share/codingclub/updatecron.sh
